# Backend/chat.service
[Unit]
Description=Ratchet Chat Service
After=default.target

[Service]
ExecStart=/usr/bin/php /var/chat/app.php

[Install]
WantedBy=default.target